# services/dorm-sports-wire.service
[Unit]
Description=Dorm Sports Wire Flask Scoreboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__PROJECT_DIR__/.venv/bin/python __PROJECT_DIR__/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
